# Phase_H/supervisord.conf
; auto-generated by Shared.deployment
; 2026-02-16T03:42:04.628973+00:00 | generated-locally
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log

[program:kalshiguard]
directory=/app
command=python Phase_A/api.py
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/tmp/kalshiguard.out.log
stderr_logfile=/tmp/kalshiguard.err.log
